# brave/Docker/brave.conf
[program:brave]
command=brave-browser
autorestart=true