# brave/Docker/xpra.conf
[program:xpra]
command=/usr/local/bin/start_xpra.sh
autostart=true
autorestart=true
user=lzc
directory=/home/lzc
environment=HOME="/home/lzc"
stopwaitsecs=3  ; 减少等待时间，默认为10秒
stopsignal=TERM